# board/trishika/dns320l/rootfs_overlay/usr/lib/systemd/system/fan.service
[Unit]
Description=Fan control

[Service]
ExecStart=/usr/local/bin/fan-daemon.sh

[Install]
WantedBy=multi-user.target
