# board/trishika/dns320l/rootfs_overlay/usr/lib/systemd/system/telnetd.service
[Unit]
Description=telnet server daemon

[Service]
ExecStart=/bin/busybox telnetd -F
Restart=always

[Install]
WantedBy=multi-user.target
